# rotator.service
[Unit]
Description=Modem Connection Rotator Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=USER_PLACEHOLDER
Group=USER_PLACEHOLDER
WorkingDirectory=INSTALL_DIR_PLACEHOLDER
ExecStart=/usr/bin/python3 INSTALL_DIR_PLACEHOLDER/rotator.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Allow the service to run networking commands
SupplementaryGroups=sudo

[Install]
WantedBy=multi-user.target